# supervisord.conf
[program:webserver]
command=gunicorn -c {{ PROJECT_DIR }}/{{ settings.GUNICORN_CONFIG }} icommons_tools.wsgi:application

[program:huey]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py run_huey
