# canvas_account_admin_tools/cron.d/canvas_account_admin_tools
* * * * * deploy /var/opt/virtualenvs/canvas_account_admin_tools/bin/python /opt/tlt/canvas_account_admin_tools/manage.py process_async_jobs >> /var/opt/tlt/logs/cron-canvas_account_admin_tools.log 2>&1
* * * * * deploy /var/opt/virtualenvs/canvas_account_admin_tools/bin/python /opt/tlt/canvas_account_admin_tools/manage.py finalize_bulk_create_jobs >> /var/opt/tlt/logs/cron-canvas_account_admin_tools.log 2>&1
